# certbot-cloudflare@.timer
[Unit]
Description=Renew certs with systemd-timer 

[Timer]
OnBootSec=300
OnUnitActiveSec=2w

[Install]
WantedBy=timers.target
